# src/consul/build/modules/build/files/etc/supervisor/conf.d/consul.conf
[program:consul]
command = /usr/local/bin/consul server -config=/etc/consul/consul.hcl
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
